# rc.local
#!/bin/sh
# Create service links on boot
ln -sf /data/service/tank-sensor /service/tank-sensor
ln -sf /data/service/bme280 /service/bme280
